# systemd/tarragon.service
[Unit]
Description=TarraGon Application Launcher
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/tarragon --daemon --port 7777
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

